# saltstack/base/php/files/etc/systemd/system/php7.2-fpm.service.d/spryker-env.conf
[Service]
ExecStartPre=/bin/mkdir -p /data/logs
ExecStartPre=/bin/mkdir -p /vagrant
ExecStartPre=/usr/bin/touch /vagrant/.vm
EnvironmentFile=/vagrant/.vm

