# saltstack/base/php/files/etc/systemd/system/php7.2-fpm.service.d/spryker-env.conf
[Service]
ExecStartPre=/bin/mkdir -p /data/logs
EnvironmentFile=/etc/spryker-vm-env
