# saltstack/base/php/files/etc/systemd/system/php7.2-fpm.service.d/spryker-env.conf
[Service]
ExecStartPre=-/bin/cp -f /vagrant/.vm /etc/spryker-vm-env
EnvironmentFile=/etc/spryker-vm-env
