# saltstack/dev/development/files/etc/cron.d/vagrant-ntpdate
PATH=/usr/bin:/usr/sbin:/bin:/sbin

1 * * * * root service ntp stop; sleep 1s; service ntp start
