# saltstack/base/nginx/files/etc/systemd/system/nginx.service.d/spryker-env.conf
[Service]
EnvironmentFile=/etc/spryker-vm-env
