# saltstack/base/system/files/etc/profile.d/fix_charset.sh
if [ "$LC_CTYPE" == "UTF-8" ]; then
  export LC_CTYPE=C
fi
