# saltstack/base/elasticsearch/files/elasticsearch_instance/usr/lib/systemd/system/elasticsearch.service.d/10-limits.conf
[Service]
LimitNPROC=16384
LimitMEMLOCK=infinity
